# device-stack/yocto-raspberrypi4/meta-cdm/recipes-cdm/cdm-enroll/files/cdm-enroll.service
[Unit]
Description=CDM Device Enrollment
Documentation=https://github.com/the78mole/complete-device-management
After=network-online.target
Wants=network-online.target
# Run only if not yet enrolled
ConditionPathExists=!/var/lib/cdm/.enrolled

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/cdm/enroll.env
ExecStart=/usr/bin/cdm-enroll.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cdm-enroll

[Install]
WantedBy=multi-user.target
